# scripts/raspberry/backup/cec.service
[Unit]
Description=MQTT2GO CEC
After=multi-user.target

[Service]
ExecStart=/bin/sh -c "/home/pi/cec.sh"
